# usr/lib/opkg/info/android-system-image-tuna.list
/system/etc/permissions/android.hardware.wifi.direct.xml
/system/etc/media_codecs.xml
/system/lib/libstagefright_amrnb_common.so
/system/xbin/mv
/system/lib/libstagefright_avc_common.so
/system/lib/libmediaplayerservice.so
/system/xbin/mktemp
/system/xbin/more
/system/xbin/inetd
/system/bin/servicemanager
/system/xbin/depmod
/system/lib/modules/ansi_cprng.ko
/system/vendor/bin/pvrsrvctl
/system/bin/printenv
/system/vendor/lib/libIMGegl.so
/system/xbin/taskset
/system/etc/audio_policy.conf
/system/addon.d/blacklist
/system/lib/libGLESv2.so
/system/bin/schedtop
/system/lib/libstagefright.so
/system/usr/keychars/qwerty.kcm
/system/lib/libext2fs.so
/system/xbin/chmod
/system/media/PFFprec_600.emd
/system/xbin/getopt
/system/xbin/od
/system/xbin/pstree
/system/lib/libcamera_compat_layer.so
/system/xbin/uname
/system/xbin/wget
/system/bin/top
/system/bin/fRom
/system/usr/share/zoneinfo/zoneinfo.idx
/system/lib/libdl.so
/system/xbin/mkfs.ext2
/system/lib/
/system/lib/libext2_uuid.so
/system/xbin/poweroff
/system/usr/
/system/lib/libstagefright_soft_amrdec.so
/system/usr/keylayout/Vendor_046d_Product_c294.kl
/system/usr/keylayout/qwerty.kl
/system/etc/bluetooth/
/system/xbin/patch
/system/bin/debuggerd
/system/xbin/crond
/system/usr/keylayout/sec_jack.kl
/system/lib/drm/
/system/etc/nfcee_access.xml
/system/xbin/sync
/system/usr/keylayout/
/system/etc/permissions/android.hardware.camera.front.xml
/system/xbin/test
/system/lib/libext2_profile.so
/system/bin/ln
/system/xbin/ps
/system/vendor/bin/
/system/xbin/base64
/system/xbin/xargs
/system/xbin/ifconfig
/system/bin/lsmod
/system/bin/ls
/system/xbin/umount
/system/lib/libandroidfw.so
/system/etc/dbus.conf
/system/lib/liblog.so
/system/vendor/lib/libinvensense_mpl.so
/system/lib/libstdc++.so
/system/bin/udptest
/system/xbin/printf
/system/xbin/uuencode
/system/lib/libstagefright_omx.so
/system/xbin/mpstat
/system/vendor/lib/hw/gps.omap4.so
/system/bin/r
/system/lib/egl/libGLES_android.so
/system/etc/dhcpcd/
/system/lib/libcorkscrew.so
/system/lib/libOMX.TI.DUCATI1.VIDEO.H264E.so
/system/xbin/free
/system/vendor/lib/libpvrANDROID_WSEGL.so
/system/etc/permissions/com.cyanogenmod.android.xml
/system/usr/keylayout/Vendor_05ac_Product_0239.kl
/system/bin/wpa_cli
/system/xbin/mke2fs
/system/lib/libis_compat_layer.so
/system/xbin/dnsd
/system/lib/libstagefright_soft_vpxdec.so
/system/etc/
/system/etc/permissions/android.hardware.camera.flash-autofocus.xml
/system/bin/mv
/system/bin/logwrapper
/system/xbin/time
/system/xbin/cal
/system/media/LMprec_508.emd
/system/lib/libwpa_client.so
/system/xbin/uptime
/system/xbin/flash_lock
/system/bin/route
/system/xbin/cat
/system/bin/dmesg
/system/usr/keychars/tuna-gpio-keypad.kcm
/system/vendor/res/images/
/system/bin/wipe
/system/lib/modules/rpmsg_client_sample.ko
/system/xbin/setconsole
/system/xbin/rm
/system/bin/drmserver
/system/usr/idc/qwerty2.idc
/system/xbin/fold
/system/bin/testid3
/system/xbin/rx
/system/usr/share/
/system/lib/libstagefright_soft_flacenc.so
/system/lib/libicui18n.so
/system/vendor/lib/hw/
/system/etc/permissions/android.hardware.location.gps.xml
/system/bin/reboot
/system/lib/hw/power.tuna.so
/system/xbin/sh
/system/lib/libspeexresampler.so
/system/xbin/tftp
/system/lib/soundfx/libbundlewrapper.so
/system/xbin/echo
/system/lib/libssl.so
/system/bin/netcfg
/system/lib/hw/audio.usb.default.so
/system/lib/libemoji.so
/system/xbin/fsync
/system/lib/libsonivox.so
/system/bin/make_ext4fs
/system/etc/wifi/wpa_supplicant.conf
/system/xbin/su
/system/xbin/crontab
/system/lib/libc_malloc_debug_leak.so
/system/bin/md5
/system/lib/libhardware_legacy.so
/system/vendor/firmware/mms144_ts_rev31.fw
/system/lib/libreference-ril.so
/system/etc/security/
/system/xbin/pidof
/system/xbin/egrep
/system/xbin/unlzma
/system/xbin/tac
/system/xbin/[
/system/lib/libstagefright_foundation.so
/system/xbin/lsmod
/system/bin/dump_image
/system/lib/libstagefrighthw.so
/system/lib/libstlport.so
/system/xbin/tr
/system/xbin/timeout
/system/xbin/tar
/system/xbin/bunzip2
/system/xbin/insmod
/system/bin/ps
/system/lib/libmedia_compat_layer.so
/system/xbin/rev
/system/xbin/grep
/system/bin/umount
/system/bin/wpa_supplicant
/system/etc/permissions/android.hardware.nfc.xml
/system/usr/icu/
/system/xbin/lsattr
/system/lib/libdrmframework.so
/system/usr/keylayout/sii9234_rcp.kl
/system/etc/NOTICE.html.gz
/system/xbin/nice
/sbin/init.android
/system/usr/keylayout/Vendor_046d_Product_c216.kl
/system/xbin/strings
/system/xbin/kill
/system/bin/uchroot
/system/xbin/usleep
/system/xbin/catv
/system/xbin/route
/system/bin/tune2fs
/system/xbin/dmesg
/system/bin/mke2fs
/sbin/adbd
/system/lib/libgui.so
/system/usr/idc/Melfas_MMSxxx_Touchscreen.idc
/system/xbin/swapoff
/system/lib/libext2_com_err.so
/system/etc/dhcpcd/dhcpcd.conf
/system/lib/libsecril-client.so
/system/xbin/sed
/system/bin/uptime
/system/xbin/lzopcat
/system/xbin/devmem
/system/lib/libutils.so
/system/xbin/vi
/system/lib/libcutils.so
/system/lib/libext4_utils.so
/system/lib/hw/power.default.so
/system/xbin/sqlite3
/system/bin/nandread
/system/bin/setconsole
/system/bin/rm
/system/bin/vmstat
/system/bin/rmmod
/system/xbin/unlzop
/system/lib/libstagefright_soft_mpeg4dec.so
/system/xbin/seq
/system/bin/linker
/system/lib/modules/
/system/lib/libcrypto.so
/system/vendor/firmware/mms144_ts_rev32.fw
/system/lib/libtinyalsa.so
/system/xbin/wc
/system/xbin/md5sum
/system/lib/hw/camera.omap4.so
/system/vendor/res/images/dock/dock.png
/system/bin/sh
/system/bin/chown
/system/xbin/xzcat
/system/etc/media_profiles.xml
/system/lib/libexpat.so
/system/xbin/fgrep
/system/etc/recovery-resource.dat
/system/lib/libnbaio.so
/system/xbin/less
/system/vendor/lib/egl/libEGL_POWERVR_SGX540_120.so
/system/lib/libexif.so
/system/xbin/killall5
/system/lib/libhardware.so
/system/bin/mount
/system/xbin/tee
/system/xbin/date
/system/etc/permissions/android.hardware.sensor.gyroscope.xml
/system/xbin/netstat
/system/lib/crtend_so.o
/system/xbin/bzip2
/system/vendor/firmware/
/system/xbin/readlink
/system/xbin/busybox
/system/etc/permissions/android.hardware.telephony.gsm.xml
/system/xbin/dhdutil
/system/xbin/head
/system/xbin/xz
/system/bin/insmod
/system/xbin/iostat
/system/etc/permissions/
/system/xbin/ntpd
/system/etc/event-log-tags
/system/xbin/stat
/system/etc/vold.fstab
/system/lib/libdsyscalls.so
/system/bin/ionice
/system/xbin/pwdx
/system/lib/libEGL.so
/system/lib/libstagefright_chromium_http.so
/system/ubuntu_stamp
/system/etc/wifi/bcmdhd.cal
/system/bin/getevent
/system/vendor/etc/audio_effects.conf
/system/bin/smc_pa_ctrl
/system/xbin/rdev
/system/vendor/lib/egl/libGLESv2_POWERVR_SGX540_120.so
/system/vendor/lib/libusc.so
/system/xbin/[[
/system/xbin/printenv
/system/xbin/renice
/system/xbin/brctl
/system/xbin/diff
/system/lib/libm.so
/system/bin/hostapd
/system/xbin/rmmod
/system/etc/permissions/android.hardware.sensor.light.xml
/system/lib/modules/scsi_wait_scan.ko
/system/lib/libthread_db.so
/system/xbin/mesg
/system/lib/libstagefright_soft_amrwbenc.so
/system/xbin/realpath
/system/lib/hw/lights.tuna.so
/system/bin/start
/system/xbin/swapon
/system/bin/ubuntu_chroot
/system/bin/mksh
/system/xbin/uniq
/system/bin/decoder
/system/xbin/chown
/system/lib/drm/libfwdlockengine.so
/system/vendor/lib/libPVRScopeServices.so
/system/lib/libstagefright_soft_vorbisdec.so
/system/etc/permissions/handheld_core_hardware.xml
/system/lib/libstagefright_soft_amrnbenc.so
/system/vendor/firmware/fw_bcmdhd_apsta.bin
/system/xbin/chattr
/system/addon.d/
/system/etc/bluetooth/main.conf
/system/lib/soundfx/libdownmix.so
/system/bin/getprop
/system/lib/soundfx/libreverbwrapper.so
/system/lib/egl/
/system/usr/keychars/sec_jack.kcm
/system/bin/stop
/system/usr/idc/autopilot-finger.idc
/system/bin/toolbox
/system/lib/libc.so
/system/xbin/mount
/system/lib/libmedia_native.so
/system/xbin/lzma
/system/lib/libskia.so
/system/lib/libmm_osal.so
/system/media/audio/ui/camera_click.ogg
/system/etc/permissions/android.hardware.usb.host.xml
/system/vendor/lib/libsrv_init.so
/system/lib/soundfx/libaudiopreprocessing.so
/system/lib/libpower.so
/system/usr/keylayout/Vendor_054c_Product_0268.kl
/system/etc/init.d/90userinit
/system/xbin/tftpd
/system/xbin/whoami
/system/vendor/lib/libsrv_um.so
/system/xbin/clear
/system/xbin/halt
/system/lib/libsync.so
/system/lib/libdrmdecrypt.so
/system/lib/libmlplatform.so
/system/lib/libI420colorconvert.so
/system/etc/permissions/android.hardware.wifi.xml
/system/bin/adb
/system/etc/audio_effects.conf
/system/lib/soundfx/libvisualizer.so
/system/addon.d/50-cm.sh
/system/bin/sendevent
/system/lib/libext2_e2p.so
/system/lib/libgabi++.so
/system/bin/sdcard
/system/lib/libbinder.so
/system/bin/sleep
/system/bin/gzip
/system/etc/permissions/android.hardware.usb.accessory.xml
/system/vendor/firmware/libpn544_fw.so
/system/lib/libstagefright_yuv.so
/system/bin/
/system/xbin/tail
/system/vendor/res/images/dock/
/system/vendor/res/
/system/xbin/unzip
/system/lib/hw/sensors.tuna.so
/system/bin/recovery
/system/lib/hw/audio_policy.default.so
/system/usr/keychars/qwerty2.kcm
/system/etc/permissions/com.tmobile.software.themes.xml
/system/bin/renice
/system/xbin/cmp
/system/xbin/dirname
/system/bin/cat
/system/usr/keychars/Virtual.kcm
/system/xbin/lzop
/system/lib/libGLES_trace.so
/system/bin/rmdir
/system/xbin/sha512sum
/system/xbin/ftpget
/system/xbin/blockdev
/system/bin/log
/system/xbin/setserial
/system/xbin/modinfo
/system/xbin/flock
/system/bin/newfs_msdos
/system/usr/keychars/
/system/bin/rild
/system/xbin/bzcat
/system/usr/keylayout/Vendor_046d_Product_c299.kl
/system/etc/permissions/android.software.sip.voip.xml
/system/lib/libext2_blkid.so
/system/xbin/adjtimex
/system/lib/libgccdemangle.so
/system/bin/erase_image
/system/xbin/ttysize
/system/vendor/firmware/fw_bcmdhd_p2p.bin
/system/vendor/firmware/bcm4330.hcd
/system/xbin/arp
/system/etc/gps.conf
/system/xbin/split
/system/lib/libOMX.TI.DUCATI1.VIDEO.DECODER.so
/system/etc/apns-conf.xml
/system/xbin/chroot
/system/lib/libstagefright_soft_h264dec.so
/system/lib/libstagefright_soft_mpeg4enc.so
/system/etc/mkshrc
/system/lib/libsparse.so
/system/xbin/ash
/system/xbin/pgrep
/system/lib/libsqlite.so
/system/bin/mediaserver
/system/xbin/groups
/system/media/bootanimation.zip
/system/lib/libc_malloc_debug_qemu.so
/system/lib/libstagefright_wfd.so
/system/xbin/sort
/system/usr/keylayout/Vendor_046d_Product_c532.kl
/system/bin/pngtest
/system/lib/libstagefright_soft_aacdec.so
/system/usr/icu/icudt48l.dat
/system/xbin/resize
/system/lib/libjpeg.so
/system/bin/watchprops
/system/vendor/lib/hw/gralloc.omap4.so
/system/xbin/sleep
/system/bin/sysinit
/system/xbin/pkill
/system/etc/permissions/android.hardware.sensor.proximity.xml
/system/xbin/mountpoint
/system/xbin/ftpput
/system/bin/installd
/system/lib/libz.so
/system/lib/libOMX_Core.so
/system/xbin/run-parts
/system/lib/libeffects.so
/system/xbin/pmap
/system/media/
/system/lib/libOMX.TI.DUCATI1.MISC.SAMPLE.so
/system/xbin/telnetd
/system/usr/share/zoneinfo/zoneinfo.dat
/system/bin/hostapd_cli
/system/lib/libstagefright_soft_rawdec.so
/system/usr/keychars/sii9234_rcp.kcm
/system/xbin/mknod
/system/bin/screencap
/system/bin/bootanimation
/system/xbin/rmdir
/system/usr/keylayout/tuna-gpio-keypad.kl
/system/vendor/lib/
/system/vendor/lib/egl/
/system/xbin/pwd
/system/lib/hw/audio.primary.tuna.so
/system/lib/libctest.so
/system/xbin/yes
/system/xbin/top
/system/xbin/cp
/system/xbin/mkswap
/system/lib/libpowermanager.so
/system/lib/libcamera_metadata.so
/system/xbin/modprobe
/system/lib/libtiutils.so
/system/xbin/basename
/system/lib/libubuntu_application_api.so
/system/bin/run-as
/system/xbin/sha1sum
/system/xbin/dc
/system/lib/libcameraservice.so
/system/xbin/dd
/system/xbin/env
/system/xbin/df
/system/media/audio/
/system/lib/libaudioutils.so
/system/bin/setprop
/system/lib/libchromium_net.so
/system/xbin/flash_unlock
/system/xbin/du
/system/lib/hw/gralloc.default.so
/system/xbin/mkfifo
/system/xbin/unix2dos
/system/xbin/dos2unix
/system/bin/lsof
/system/xbin/traceroute
/system/lib/libinvensense_hal.so
/system/bin/mkdir
/system/usr/keylayout/AVRCP.kl
/system/xbin/sysctl
/system/xbin/unexpand
/system/xbin/awk
/system/xbin/reset
/system/xbin/add-property-tag
/system/xbin/fuser
/system/bin/sync
/system/usr/idc/
/system/vendor/etc/
/system/xbin/ed
/system/bin/bugmailer.sh
/system/bin/netstat
/system/xbin/losetup
/system/xbin/freeramdisk
/system/bin/ifconfig
/system/vendor/firmware/fw_bcmdhd.bin
/system/xbin/nanddump
/system/lib/libvorbisidec.so
/system/bin/flash_image
/system/vendor/etc/sirfgps.conf
/system/xbin/hexdump
/system/lib/hw/nfc.tuna.so
/system/lib/libOMX.TI.DUCATI1.VIDEO.CAMERA.so
/system/lib/libsurfaceflinger_ddmconnection.so
/system/lib/libstagefright_enc_common.so
/system/bin/logcat
/system/bin/backuptool.functions
/system/bin/backuptool.sh
/system/etc/spn-conf.xml
/system/lib/libGLESv1_CM.so
/system/xbin/chgrp
/system/lib/libril.so
/system/etc/init.d/00banner
/system/bin/handle_compcache
/system/bin/surfaceflinger
/system/xbin/man
/system/usr/share/zoneinfo/zoneinfo.version
/system/etc/permissions/com.android.nfc_extras.xml
/system/xbin/expr
/system/etc/init.d/
/system/xbin/expand
/system/lib/libwebrtc_audio_preprocessing.so
/system/bin/compcache
/system/xbin/cut
/system/xbin/nandwrite
/system/bin/restart
/system/usr/keylayout/Vendor_22b8_Product_093d.kl
/system/xbin/ping
/system/bin/send_bug
/system/xbin/sum
/system/vendor/lib/hw/hwcomposer.omap4.so
/system/vendor/lib/libpvr2d.so
/system/etc/security/otacerts.zip
/system/vendor/lib/libsec-ril.so
/system/xbin/bbconfig
/system/bin/cp
/system/media/audio/ui/
/system/lib/libsysutils.so
/system/xbin/lzcat
/system/etc/permissions/android.hardware.audio.low_latency.xml
/system/vendor/lib/egl/libGLESv1_CM_POWERVR_SGX540_120.so
/system/xbin/false
/system/lib/libmedia.so
/system/xbin/uudecode
/system/lib/libion_ti.so
/system/lib/hw/local_time.default.so
/system/bin/dd
/system/bin/setup_fs
/system/vendor/firmware/hdcp.keys
/system/bin/df
/system/xbin/mkdir
/system/lib/libETC1.so
/system/lib/libinput.so
/system/etc/permissions/com.cyanogenmod.nfc.enhanced.xml
/default.prop
/system/xbin/flashcp
/system/xbin/find
/system/bin/notify
/system/bin/du
/system/lib/libnetutils.so
/system/etc/permissions/android.hardware.touchscreen.multitouch.jazzhand.xml
/system/xbin/lsusb
/system/xbin/cpio
/system/bin/touch
/system/xbin/id
/system/lib/crtbegin_so.o
/system/lib/hw/
/system/lib/libpixelflinger.so
/system/bin/grep
/system/xbin/ip
/system/lib/libstagefright_soft_h264enc.so
/system/xbin/gzip
/system/xbin/nslookup
/system/lib/libstagefright_soft_mp3dec.so
/system/etc/hosts
/system/lib/soundfx/
/system/lib/libsurfaceflinger.so
/system/etc/CHANGELOG-CM.txt
/system/xbin/
/system/lib/libtf_crypto_sst.so
/system/usr/keylayout/Generic.kl
/system/usr/share/zoneinfo/
/system/xbin/length
/system/bin/ioctl
/system/xbin/killall
/system/xbin/tune2fs
/system/lib/libstagefright_soft_g711dec.so
/system/bin/kill
/system/vendor/firmware/ducati-m3.bin
/system/etc/permissions/com.nxp.mifare.xml
/system/lib/libstagefright_soft_aacenc.so
/system/lib/egl/egl.cfg
/system/lib/libOMX.TI.DUCATI1.VIDEO.DECODER.secure.so
/system/etc/init.goldfish.sh
/system/bin/cmp
/system/lib/libui.so
/system/xbin/comm
/system/xbin/telnet
/system/bin/smd
/system/xbin/fdisk
/system/xbin/which
/system/bin/chmod
/system/xbin/unxz
/system/lib/modules/rpmsg_server_sample.ko
/system/vendor/lib/libglslcompiler.so
/system/xbin/install
/system/lib/libsensorservice.so
/vendor
/system/lib/hw/audio.primary.default.so
/system/vendor/
/system/bin/applypatch
/system/lib/libcommon_time_client.so
/system/lib/libdomx.so
/system/lib/libsf_compat_layer.so
/system/lib/libicuuc.so
/system/bin/e2fsck
/system/bin/ubuntuappmanager
/system/etc/wifi/
/system/bin/wfd
/system/bin/iftop
/system/xbin/setsid
/system/usr/keylayout/Vendor_045e_Product_028e.kl
/system/xbin/sha256sum
/system/xbin/uncompress
/system/xbin/zcat
/system/bin/sensorservice
/system/lib/libcamera_client.so
/system/bin/hd
/system/lib/libusbhost.so
/system/xbin/watch
/system/usr/keychars/Generic.kcm
/system/xbin/ln
/system/lib/libOMX.TI.DUCATI1.VIDEO.MPEG4E.so
/system/xbin/gunzip
/system/xbin/ls
/system/build.prop
/system/
/system/xbin/true
/system/xbin/touch
/system/bin/radiooptions
/system/xbin/blkid
/system/lib/libaudioflinger.so
/system/bin/date
/system/etc/permissions/android.hardware.sensor.barometer.xml
/system/xbin/nohup
/system/xbin/mkfs.vfat
/system/usr/idc/qwerty.idc
/system/xbin/check-lost+found
/system/xbin/stty
/system/lib/libmllite.so
/system/bin/id
