# usr/lib/opkg/info/token-generator.list
/etc/rcS.d/
/etc/rcS.d/S65token-generator.sh
/etc/init.d/token-generator.sh
