# usr/lib/opkg/info/update-rc.d-lic.list
/usr/share/licenses/update-rc.d/update-rc.d
/usr/share/licenses/update-rc.d/
/usr/share/licenses/update-rc.d/generic_GPLv2
