# usr/lib/opkg/info/upstart-sysvcompat.list
/etc/event.d/control-alt-delete
/etc/default/rcS
/sbin/halt.upstart
/etc/init.d/rc
/sbin/poweroff.upstart
/etc/event.d/rcS-sulogin
/etc/event.d/rc-default
/sbin/shutdown.upstart
/sbin/reboot.upstart
/etc/init.d/rcS
/sbin/runlevel
/etc/event.d/sulogin
/etc/event.d/rc0
/etc/event.d/rc1
/etc/event.d/rc2
/etc/event.d/rc3
/etc/event.d/rc4
/etc/event.d/rc5
/etc/event.d/rc6
/etc/event.d/rcS
